# deploy/us.quorum.qherlock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>us.quorum.qherlock</string>

    <!-- EDIT: replace __REPO__ with the absolute repo path. WorkingDirectory is
         load-bearing: qherlock reads .env, data/, and runs/ relative to it. -->
    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <!-- zsh -lc loads the login profile so python3/uv/tsh are on PATH.
         uv is module-installed on this host: python3 -m uv. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>exec python3 -m uv run qherlock patrol --scope all</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>7</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <!-- EDIT: replace __HOME__ (launchd does not expand ~).
         mkdir -p ~/Library/Logs/qherlock first. -->
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/qherlock/patrol.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/qherlock/patrol.err.log</string>
</dict>
</plist>
